# deploy/gateway-torpool.service
[Unit]
Description=Gateway TorPool Manager
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=gateway
Group=gateway
ExecStart=/opt/gateway/gateway-torpool -config /etc/gateway/config.yaml
Restart=always
RestartSec=5
LimitNOFILE=65536

# Directories — systemd creates/cleans these automatically
RuntimeDirectory=gateway
RuntimeDirectoryMode=0750
StateDirectory=gateway/tor
StateDirectoryMode=0700
LogsDirectory=gateway
LogsDirectoryMode=0750

# Security hardening
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/var/lib/gateway /var/log/gateway /run/gateway /tmp

[Install]
WantedBy=multi-user.target
